# ansible/roles/microservice/templates/metrics.service.j2
[Unit]
Description=Prometheus metrics HTTP server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/metrics_app.py
Restart=always

[Install]
WantedBy=multi-user.target
